# dotfiles/arch/sway-session.target
# vim:set filetype=systemd:

[Unit]
Description=Sway session
